# daily-briefing/com.frank.dailybriefing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.frank.dailybriefing</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>BRIEFING_SCRIPT_PATH</string>
    </array>

    <!-- Fire at 6:00 AM every day (local time) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>BRIEFING_LOG_PATH</string>
    <key>StandardErrorPath</key>
    <string>BRIEFING_LOG_PATH</string>

    <!-- Don't run immediately on load, only at scheduled time -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
